# sentinelzero.service
[Unit]
Description=SentinelZero Network Security Scanner
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/home/sentinel/SentinelZero/backend
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=FLASK_ENV=production
Environment=SENTINEL_BIND_HOST=0.0.0.0
Environment=SENTINEL_BIND_PORT=5000
Environment=SENTINEL_ALLOWED_ORIGINS=http://localhost:5000,http://127.0.0.1:5000,http://sentinelzero.ops.prox,http://sentinelzero.ops.prox:5000,https://sentinelzero.ops.prox,http://172.16.0.198:5000
Environment=SCAN_TIMEOUT_SECONDS=1800
Environment=SCAN_TIMEOUT_HOME_SECONDS=5400
Environment=SCAN_TIMEOUT_IOT_HOME_SECONDS=7200
Environment=SCAN_TIMEOUT_VULN_SECONDS=10800
Environment=SCAN_XML_PARSE_RETRIES=3
Environment=SCAN_XML_PARSE_RETRY_DELAY_SECONDS=1
ExecStart=/home/sentinel/SentinelZero/backend/run-app.sh
ExecStop=/bin/bash -c "pkill -f nmap; pkill -f 'python.*app.py'"
Restart=on-failure
RestartSec=5
TimeoutStartSec=45
TimeoutStopSec=20
KillMode=mixed
StandardOutput=journal
StandardError=journal

# Security settings - relaxed for network scanning
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/sentinel/SentinelZero
# Allow raw socket access for nmap
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
