# systemd/sentinelzero-backend.service
[Unit]
Description=SentinelZero Backend
After=network.target

[Service]
Type=simple
User=sentinelzero
Group=sentinelzero
WorkingDirectory=/home/sentinel/SentinelZero/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/sentinel/SentinelZero/backend/.venv/bin/python app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
